# trading-bot/com.user.tradingbot.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!--
  macOS launchd schedule for the trading bot.

  This runs the bot ONCE per weekday at 10:00 (after the market open settles).
  A slow, scheduled cadence is intentional: it respects the Pattern Day Trader
  rule and avoids spread/slippage churn that destroys small accounts.

  SETUP:
    1. Edit the two paths below (ProgramArguments and WorkingDirectory) to your
       real locations. Use absolute paths.
    2. Copy to ~/Library/LaunchAgents/com.user.tradingbot.plist
    3. Load it:    launchctl load ~/Library/LaunchAgents/com.user.tradingbot.plist
    4. Unload it:  launchctl unload ~/Library/LaunchAgents/com.user.tradingbot.plist

  NOTE: your Mac must be awake at the scheduled time for the run to fire. If it
  is asleep, launchd runs it at the next wake. Keep the machine powered/awake,
  or run the bot on an always-on host.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.tradingbot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/Dylan/.local/bin/python3</string>
        <string>/Users/Dylan/path/to/milestoneiq-app/trading-bot/run.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/Dylan/path/to/milestoneiq-app/trading-bot</string>

    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>10</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>10</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>10</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>10</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>10</integer><key>Minute</key><integer>0</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/Dylan/path/to/milestoneiq-app/trading-bot/data/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/Dylan/path/to/milestoneiq-app/trading-bot/data/launchd.err.log</string>
</dict>
</plist>
